# install/roles/file-service/templates/file-service-cleaner/crontab/root.j2
*/15 * * * * /usr/bin/bash /clear_deleted_dir.sh /file-service/data
